# deployment/systemd/auramail-web.service
[Unit]
Description=AuraMail Flask Web Application
After=network.target redis.service postgresql.service

[Service]
Type=notify
User=auramail
Group=auramail
WorkingDirectory=/opt/auramail
Environment="PATH=/opt/auramail/venv/bin"
EnvironmentFile=/etc/auramail/auramail.env
ExecStart=/opt/auramail/venv/bin/gunicorn \
    --workers 4 \
    --bind 127.0.0.1:5000 \
    --timeout 120 \
    --access-logfile /var/log/auramail/access.log \
    --error-logfile /var/log/auramail/error.log \
    server:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target


